# deploy/crontab
# Cardnews 자동 실행 — KST 08:00 (UTC 23:00), KST 13:00 (UTC 04:00)
0 23 * * * cd /app && /usr/local/bin/python main.py >> /app/logs/cron.log 2>&1
0 4 * * * cd /app && /usr/local/bin/python main.py >> /app/logs/cron.log 2>&1
